# assets/config/supervisor/6-socat-cdp.conf
[program:socat-cdp]
command=socat TCP-LISTEN:9222,fork TCP:localhost:9223
autorestart=true
stderr_logfile=/var/log/socat-cdp.err.log
priority=6
